# templates/default/etc/systemd/system/pgbouncer.service.erb
[Unit]
Description=A lightweight connection pooler for PostgreSQL
After=network.target

[Service]
User=pgbouncer
ExecStart=/usr/sbin/pgbouncer /etc/pgbouncer/pgbouncer-<%= @db_alias %>.ini
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

